# templates/default/tfc-agent.service.erb
[Unit]
Description=HashiCorp Terraform Cloud Agent

[Service]
Type=simple
ExecStart=/usr/bin/tfc-agent -token <%= node['tfc_agent']['token'] %> 

[Install]
WantedBy=multi-user.target